# crontab
#run python script every 30th minute
*/30 * * * * python /code/pagespeed10x/pagespeed10x.py > /proc/1/fd/1 2>/proc/1/fd/2
